# Docs/Published/Deployment/launchd/com.tldw.worker.chatbooks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Generated by Helper_Scripts/Deployment/generate_sidecar_files.py. Edit Docs/Deployment/sidecar_workers_manifest.json and Docs/Deployment/Sidecar_Workers.template.md. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tldw.worker.chatbooks</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/tldw_server/.venv/bin/python</string>
    <string>-m</string>
    <string>tldw_Server_API.app.services.core_jobs_worker</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/opt/tldw_server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>UserName</key>
  <string>tldw</string>
  <key>GroupName</key>
  <string>tldw</string>
  <!-- Ensure /opt/tldw_server/logs/launchd exists and is writable by tldw:tldw. -->
  <key>StandardOutPath</key>
  <string>/opt/tldw_server/logs/launchd/chatbooks.log</string>
  <key>StandardErrorPath</key>
  <string>/opt/tldw_server/logs/launchd/chatbooks.err.log</string>
</dict>
</plist>
