# Docs/Published/Deployment/systemd/tldw-worker-reading-import.timer
# Generated by Helper_Scripts/Deployment/generate_sidecar_files.py. Edit Docs/Deployment/sidecar_workers_manifest.json and Docs/Deployment/Sidecar_Workers.template.md.
[Unit]
Description=Ensure tldw worker (reading-import) is running

[Timer]
OnBootSec=30s
OnStartupSec=30s
Unit=tldw-worker-reading-import.service
Persistent=true

[Install]
WantedBy=timers.target
