# Docs/Published/Deployment/systemd/tldw-worker-chatbooks.timer
# Generated by Helper_Scripts/Deployment/generate_sidecar_files.py. Edit Docs/Deployment/sidecar_workers_manifest.json and Docs/Deployment/Sidecar_Workers.template.md.
[Unit]
Description=Ensure tldw worker (chatbooks) is running

[Timer]
OnBootSec=30s
OnUnitActiveSec=5m
Unit=tldw-worker-chatbooks.service
Persistent=true

[Install]
WantedBy=timers.target
