# Docs/site/Deployment/launchd/com.tldw.worker.prompt-studio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Generated by Helper_Scripts/Deployment/generate_sidecar_files.py. Edit Docs/Deployment/sidecar_workers_manifest.json and Docs/Deployment/Sidecar_Workers.template.md. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tldw.worker.prompt-studio</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/tldw_server/.venv/bin/python</string>
    <string>-m</string>
    <string>tldw_Server_API.app.core.Prompt_Management.prompt_studio.services.jobs_worker</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/opt/tldw_server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/opt/tldw_server/logs/launchd/prompt-studio.log</string>
  <key>StandardErrorPath</key>
  <string>/opt/tldw_server/logs/launchd/prompt-studio.err.log</string>
</dict>
</plist>
